# mgz/supervisor/app.conf
[program:flaskapp]
command=/usr/bin/runapp.sh
autostart=true
autorestart=true